# devel/etc/catkin/profile.d/25.turtlebot.sh
/home/martinelli/WidowX_XM430_Robot_Turret/catkin_ws/devel/.private/turtlebot_bringup/etc/catkin/profile.d/25.turtlebot.sh