# devel/etc/catkin/profile.d/25.turtlebot-gazebo.sh
/home/martinelli/WidowX_XM430_Robot_Turret/catkin_ws/devel/.private/turtlebot_gazebo/etc/catkin/profile.d/25.turtlebot-gazebo.sh